# ansible/roles/server/templates/prefect-worker.service.j2
[Unit]
Description=Prefect worker (host-pool) for host jobs
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory={{ prefect_flow_code_host_path }}
ExecStart={{ prefect_venv_path }}/bin/prefect worker start --pool host-pool
Environment=PREFECT_API_URL=http://127.0.0.1:{{ prefect_server_port }}/api
Environment=REGISTRY_URL=registry.{{ infrastructure_secrets.base_domain }}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
